# rootfs/etc/systemd/system/postiz-db-setup.service
[Unit]
Description=Create Postiz and Temporal databases
After=postgresql.service
Requires=postgresql.service

[Service]
Type=oneshot
User=postgres
ExecStart=/usr/local/bin/postiz-db-setup.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
